# flatpak-update.service
[Unit]
Description=flatpaksをアップデート

[Service]
ExecStart=flatpak update --noninteractive
