# flatpak-update.timer
[Unit]
Description=flatpaksをアップデート

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target